# packages/vault-sync/service-units/systemd/wiki-fetch.service
[Unit]
Description=vault-sync fetch (git fetch + notify)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=@SCRIPT_DIR@/wiki-fetch-notify.sh
RuntimeMaxSec=120
Nice=10
