# packages/vault-sync/service-units/systemd/wiki-fuse-refresh.service
[Unit]
Description=vault-sync fuse refresh (rclone rc vfs/refresh)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=HOME=@HOME@
ExecStart=@SCRIPT_DIR@/wiki-fuse-refresh.sh
RuntimeMaxSec=120
Nice=10
