# packages/vault-sync/service-units/launchd/com.karlchow.wiki-fetch.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- TEMPLATE — placeholders replaced by vault-sync-install:
     @SCRIPT_DIR@      → $HOME/.local/share/vault-sync/bin
     @LOG_DIR@         → $HOME/Library/Logs
     @LAUNCHD_PATH@    → installer PATH containing node + fallback bins
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.karlchow.wiki-fetch</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>@SCRIPT_DIR@/wiki-fetch-notify.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>@LOG_DIR@/wiki-fetch.log</string>
    <key>StandardErrorPath</key>
    <string>@LOG_DIR@/wiki-fetch.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>@LAUNCHD_PATH@</string>
        <key>HOME</key>
        <string>/Users/karlchow</string>
        <key>WIKI_FETCH_PULL_ON_DELTA</key>
        <string>1</string>
    </dict>
</dict>
</plist>
